# systemd_services/activitypub.socket
[Unit]
Description=activitypub gunicorn socket

[Socket]
ListenStream=/run/activitypub.sock

[Install]
WantedBy=sockets.target
